# com.system76.FirmwareManager.target
[Unit]
Description=Ensures that the fwupd service has been started
Requires=network-online.target
Wants=fwupd.service
