# debian/firmware-manager-notify.postinst
#!/bin/sh

systemctl --global enable com.system76.FirmwareManager.Notify.timer

#DEBHELPER#

exit 0