# debian/firmware-manager-notify.postrm
#!/bin/sh
set -x
systemctl --global disable com.system76.FirmwareManager.Notify.timer

#DEBHELPER#

exit 0